# exe-gateway.service
# exe-gateway — systemd service unit
#
# Install:
#   cp exe-gateway.service /etc/systemd/system/
#   systemctl daemon-reload
#   systemctl enable --now exe-gateway
#
# Runtime configuration:
#   /etc/exe-gateway/exe-gateway.env   - secrets + env overrides
#   /home/exe/.exe-os/gateway.json     - adapter/runtime config
#
# Logs:
#   journalctl -u exe-gateway -f

[Unit]
Description=exe-gateway multi-platform messaging gateway
Documentation=https://github.com/AskExe/exe-gateway
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=exe
Group=exe
WorkingDirectory=/opt/exe-gateway

Environment=NODE_ENV=production
EnvironmentFile=/etc/exe-gateway/exe-gateway.env
ExecStart=/usr/bin/node /opt/exe-gateway/dist/bin/exe-gateway.js

Restart=always
RestartSec=5
StartLimitIntervalSec=300
StartLimitBurst=10

StandardOutput=journal
StandardError=journal
SyslogIdentifier=exe-gateway

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/exe/.exe-os
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

LimitNOFILE=65536
MemoryMax=4G

[Install]
WantedBy=multi-user.target
